# phosphor-ipmi-net@.socket
[Unit]
ConditionPathExists=/sys/class/net/%i

[Socket]
ListenDatagram=623
BindToDevice=%i

[Install]
WantedBy=sockets.target

